# deploy/systemd-user/lark-memory-compute.service
[Unit]
Description=LarkMemoryCore Compute Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/lark-memory-core
Environment=HOME=%h
EnvironmentFile=%h/lark-memory-core/.env
ExecStart=/bin/bash -lc 'exec ./build/bin/compute_server'
Restart=always
RestartSec=5
TimeoutStopSec=70
KillMode=mixed
StandardOutput=append:%h/lark-memory-core/.run/systemd-compute.log
StandardError=append:%h/lark-memory-core/.run/systemd-compute.log

[Install]
WantedBy=default.target
