# deploy/systemd-user/ruyi-api.service
[Unit]
Description=Ruyi Serving API Server
After=network-online.target ruyi-compute.service
Wants=network-online.target
Requires=ruyi-compute.service

[Service]
Type=simple
WorkingDirectory=%h/ruyi-serving
Environment=HOME=%h
EnvironmentFile=%h/ruyi-serving/.env
ExecStart=/bin/bash -lc 'exec python3 -m api_server.main'
Restart=always
RestartSec=5
TimeoutStopSec=40
KillMode=mixed
StandardOutput=append:%h/ruyi-serving/.run/systemd-api.log
StandardError=append:%h/ruyi-serving/.run/systemd-api.log

[Install]
WantedBy=default.target
